# files/etc/systemd/system/gnome-remote-desktop.service.d/debug.conf.example
[Service]
Environment=G_MESSAGES_DEBUG=all
Environment=GNOME_REMOTE_DESKTOP_DEBUG=1